# contrib/systemd/vpn-failover-user.timer
[Unit]
Description=Run sing-box VPN failover probe every 5 minutes (per-user)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30
Persistent=true
Unit=vpn-failover-user.service

[Install]
WantedBy=timers.target
