# contrib/tproxy/systemd/sing-box-watchdog.timer
[Unit]
Description=sing-box TPROXY watchdog timer (every 1 minute)

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=sing-box-watchdog.service

[Install]
WantedBy=timers.target